# srcs/mysql/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=sh -c "mysqld_safe --skip-networking=off && kill -s SIGTERM $(cat supervisord.pid)"
startsecs=0
redirect_stderr=true
